# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:backend]
command=/usr/bin/node /app/backend/server.js
autostart=true
autorestart=true
